# babeld/install/etc/systemd/system/orbis-network.service
[Unit]
Description=ORBIS Mesh Network Base Setup (bridge + mesh + IPs)
After=network.target orbis-ssh-ip.service
Wants=network.target orbis-ssh-ip.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/orbis/network/orbis-network.sh start
ExecStop=/opt/orbis/network/orbis-network.sh stop
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
