# babeld/install/etc/systemd/system/orbis-ui.service
[Unit]
Description=ORBIS Web Interface (Flask via gunicorn)
After=network.target orbis-network.service
Wants=orbis-network.service

[Service]
Type=simple
WorkingDirectory=/opt/orbis/interface
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/orbis/interface/run-ui.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
