# batman - Backup/install/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
[Service]
# Clear inherited ExecStart and set a faster behavior: accept any interface as online
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --any --timeout=30
TimeoutStartSec=35
